# taifish.service
[Unit]
Description=泰摸鱼吧 - DBA数据库管理Web应用
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/TaifishingV4
Environment=PATH=/path/to/TaifishingV4/venv/bin
Environment=FLASK_APP=app.py
Environment=FLASK_ENV=production
Environment=REDIS_URL=redis://localhost:6379/0
ExecStart=/path/to/TaifishingV4/venv/bin/python start_production.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taifish

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/TaifishingV4/userdata
ReadWritePaths=/path/to/TaifishingV4/logs

[Install]
WantedBy=multi-user.target
